# tools/crawler/mcprs-crawl.service
[Unit]
Description=Crawl MCPRS packages

[Service]
Type=oneshot
ExecStart=bash ./crawl.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
